# dotfiles/profile.d/0-common
# vi: ft=sh

# Common env vars and settings
export EDITOR=vim

