# dotfiles/profile.d/dot-local
# vi: ft=sh

# Dot-local path

if [ -d "$HOME/.local/bin" ]; then
    PATH="$PATH:$HOME/.local/bin"
    export PATH
fi
